# device-ui-flutter/packaging/meetingbox-device-bridge.service
[Unit]
Description=MeetingBox device-services bridge (local hardware API for Flutter UI)
After=network.target
Wants=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/meetingbox/device-ui-flutter.env
Environment=DEVICE_BRIDGE_HOST=127.0.0.1
Environment=DEVICE_BRIDGE_PORT=8765
ExecStart=/usr/lib/meetingbox/device-services/venv/bin/python -m uvicorn bridge:app --app-dir /usr/lib/meetingbox/device-services/src --host ${DEVICE_BRIDGE_HOST} --port ${DEVICE_BRIDGE_PORT}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
